# mscs
#!/bin/sh
### BEGIN INIT INFO
# Provides:       mscs
# Required-Start: $remote_fs $syslog
# Required-Stop:  $remote_fs $syslog
# Default-Start:  2 3 4 5
# Default-Stop:   0 1 6
# chkconfig:      345 50 50
# Description:    Minecraft Server Control Script
### END INIT INFO

. /etc/default/mscs
USER_NAME=${USER_NAME:-minecraft}

sudo -u $USER_NAME msctl -c /etc/default/mscs "$@"
